# supervisor_conf/development/djangoapp.conf
; ==================================
;  Django app supervisor
; ==================================

[program:django_app]
; Directory should become before command
directory=/myapp/cyto_soft_mgr/

; user=nobody
numprocs=1
stdout_logfile=/var/log/django_app/logs.log
stderr_logfile=/var/log/django_app/logs.log
autostart=true
autorestart=true
startsecs=10

command=python3 manage.py runserver 0.0.0.0:%(ENV_DJANGO_PORT)s

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
; stopwaitsecs = 900
stopwaitsecs = 10

; Causes supervisor to send the termination signal (SIGTERM) to the whole process group.
stopasgroup=true

; Set Celery priority higher than default (999)
; so, if rabbitmq is supervised, it will start first.
priority=999
